# templates/default/revealcloud.initd.erb
#!/sbin/runscript
# Author: CopperEgg <support@copperegg.com>
#

# install defaults
INSTALL_ROOT="/usr/local/revealcloud"
INSTALL_BIN="$INSTALL_ROOT/revealcloud"
RUNDIR="$INSTALL_ROOT/run"
PIDFILE="$RUNDIR/revealcloud.pid"
EXTRA_ARGS=""

API_HOST="api.copperegg.com"
APP_HOST="app.copperegg.com"

RCUSER=revealcloud
RCGROUP=${RCUSER}

RC_APIKEY='<%= @apikey %>'
RC_PROXY='<%= @proxy %>'
RC_OOM_PROTECT='<%= @oom_protect %>'
RC_UUID='<%= @use_fqdn %>'
RC_TAG='<%= @tags %>'
RC_LABEL='<%= @label %>'
RC_FQDN='<%= @rc_fqdn %>'

if [ ! -z "$RC_PROXY" ]; then
    PROXY_SETTING="-p "${RC_PROXY}
    CURL_PROXY_SETTING="-x "${RC_PROXY}
fi

if [ ! -z "$RC_LABEL" ]; then
    LABEL_SETTING="-l "${RC_LABEL}
fi

if [ ! -z "$RC_TAG" ]; then
    TAG_SETTING=""
    for tag in `echo ${RC_TAG} | sed -e 's/,/ /g'`; do
        TAG_SETTING="${TAG_SETTING} -t "${tag}
    done
fi

if [ ! -z "$RC_UUID" ]; then
    # If RC_UUID is set to 0/false/no/null, do nothing.
    # Else pass the fqdn to revealcloud using -U
    case "$RC_UUID" in
        0|f*|n*)
            UUID_SETTING=""
        ;;
        *)
            UUID_SETTING="-U "${RC_FQDN}
        ;;
    esac
fi

if [ ! -z "$RC_OOM_PROTECT" ]; then
    # If OOM_PROTECT is set to 0/false/no/null, do nothing.
    # Else enable oom protection with -O
    case "$RC_OOM_PROTECT" in
        0|f*|n*)
            OOM_SETTING=""
        ;;
        *)
            OOM_SETTING="-O "
        ;;
    esac
fi

DESC="analytics collection for CopperEgg revealcloud"
NAME=revealcloud
DAEMON=$INSTALL_BIN
DAEMON_ARGS="-d -x -a ${API_HOST} -k ${RC_APIKEY} ${PROXY_SETTING} ${OOM_SETTING} ${UUID_SETTING} ${TAG_SETTING} ${LABEL_SETTING} -P ${RUNDIR} ${EXTRA_ARGS}"


opts="start stop restart"

depend() {
    use net
}

start() {
    ebegin "Starting revealcloud"
    start-stop-daemon --start --quiet \
    			--pidfile ${PIDFILE} \
    			--chuid ${RCUSER}:${RCGROUP} \
    			--exec ${DAEMON} -- \
    			${DAEMON_ARGS}
	eend 0
}

stop() {
    ebegin "Stopping revealcloud"
    start-stop-daemon --stop --quiet \
    			--pidfile ${PIDFILE} \
				--chuid ${RCUSER}:${RCGROUP} \
    			--retry=TERM/3/KILL/5 \
    			 --name ${NAME}
    local rslt=0
    rslt=$?
    if [ "$rslt" -ne  0 ]; then
   		# Wait for children to finish
    	start-stop-daemon --stop --quiet \
    				--retry=0/30/KILL/5 \
    				--exec ${DAEMON}
    fi
    # Many daemons don't delete their pidfiles when they exit.
    rm -f ${PIDFILE}
	eend 0
}

restart() {
    ebegin "Restarting revealcloud"
    svc_stop
    svc_start
    eend $? "Failed to restart revealcloud"
}
